# butrmanager/butrmanager.conf
[program:butrmanager]
command=python3 butrserver.py
directory=/path/to/butrmanager
autostart=True
autorestart=True
startsecs=5
startretries=20
redirect_stderr=True
stdout_logfile=AUTO
